# deploy/systemd/bigip-ts-validator.service
[Unit]
Description=BIG-IP Telemetry Streaming validator (API + static Web UI)
Documentation=https://github.com/gregcoward/bigip_ts_validator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bigip-ts
Group=bigip-ts
WorkingDirectory=/opt/bigip-ts-validator
Environment=PATH=/opt/bigip-ts-validator/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/bigip-ts-validator/.venv/bin/uvicorn server.app:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=true

[Install]
WantedBy=multi-user.target
